# crontab
*/5 * * * * root /bin/bash -c /root/script.sh >> /var/log/cron.log 2>&1
